# v8/build/linux/debian_sid_i386-sysroot/lib/systemd/system/speech-dispatcherd.service
[Unit]
Description=Speech-Dispatcher, common interface to speech synthesizers

[Service]
Type=forking
User=speech-dispatcher
PIDFile=/var/run/speech-dispatcher/speech-dispatcher.pid
ExecStart=/usr/bin/speech-dispatcher -d --pid-file /var/run/speech-dispatcher/speech-dispatcher.pid
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
Alias=speech-dispatcher.service
